# ArchLinux/systemd/unattended-update.timer
[Unit]
Description=Unattended Update

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3hours

[Install]
WantedBy=multi-user.target
